# config/blixa/smidig2009.smf.xml.erb
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='network/mongrel/smidig2010' type='service' version='0'>
    <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <dependency name='net' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/loopback'/>
    </dependency>
    <exec_method name='start' type='method' exec='/opt/csw/bin/mongrel_rails cluster::start --config tmp/mongrel_cluster.yml' timeout_seconds='60'>
      <method_context/>
    </exec_method>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'>
      <method_context/>
    </exec_method>
    <instance name='<%= environment %>' enabled='true'>
      <method_context project=':default' resource_pool=':default' working_directory='/u/apps/smidig2010/<%= environment %>/smidig2010'>
        <method_credential group='deploy' limit_privileges=':default' privileges=':default' supp_groups=':default' user='deploy'/>
        <method_environment>
          <envvar name='PATH' value='/local/bin:/usr/bin:/bin:/opt/csw/bin'/>
        </method_environment>
      </method_context>
      <property_group name='general' type='framework'>
        <property name='action_authorization' type='astring'/>
      </property_group>
    </instance>
  </service>
</service_bundle>
